# provisioning/openwrt-overlay/etc/uci-defaults/97-easymanet-management-lan
#!/bin/sh

if [ -x /etc/init.d/easymanet-management-lan ]; then
    /etc/init.d/easymanet-management-lan enable
fi

exit 0
